# deploy/systemd/mofa-fm-django.service
[Unit]
Description=MoFA-FM Django Application
After=network.target postgresql.service redis.service

[Service]
Type=notify
User=www-data
Group=www-data
WorkingDirectory=/opt/mofa-fm/backend
Environment="PATH=/opt/mofa-fm/backend/venv/bin"
ExecStart=/opt/mofa-fm/backend/venv/bin/gunicorn \
    --workers 4 \
    --bind unix:/run/mofa-fm.sock \
    --timeout 120 \
    --access-logfile /opt/mofa-fm/backend/logs/gunicorn-access.log \
    --error-logfile /opt/mofa-fm/backend/logs/gunicorn-error.log \
    --log-level info \
    config.wsgi:application

# 自动重启
Restart=always
RestartSec=5

# 性能优化
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
